# rc_scripts/mlmmjadmin.openbsd
#!/bin/sh

daemon='python /opt/mlmmj-admin/mlmmj-admin.py'

. /etc/rc.d/rc.subr

rc_reload=NO
rc_cmd $1

rc_stop() {
    kill `cat /var/run/mlmmj-admin.pid`
}
